# shellinabox.conf
[program:shellinabox]
command=shellinaboxd -s /:LOGIN -t -p 9002
